# ficheirosPostos/twscreen/TWDrvStartup.ORIG
#!/bin/bash
#
# Copyright 2003-2010 3M. All rights reserved.
#
# TWDrvStartup
#
# Author: 3M Touch Systems, Inc., Methuen, MA 01488 USA
#
# processname: TWDrvStartup
# chkconfig: 5 20 05
# description: Installs X input driver for the MicroTouch MT 7 Software system
#
### BEGIN INIT INFO
# Provides: TwDriver
# Required-Start:
# Required-Stop:
# Default-Start: 2 5
# Default-Stop: 0 1 3 4 6
# Description: Start the MT 7 touch screen driver
### END INIT INFO

#
# Add startup options for TwDriver here
#
TWDRV_OPTIONS=""

#
# Determine how we can emit status, giving preference to the
# LSB process
#
if test -f /lib/lsb/init-functions
then
   . /lib/lsb/init-functions
   OUTPUT_TYPE="3"
elif test -f %LSBDIR%/init-functions
then
   . %LSBDIR%/init-functions
   OUTPUT_TYPE="3"
elif test -f /etc/init.d/functions
then
   . /etc/init.d/functions
   OUTPUT_TYPE="2"
elif test -f %INITDIR%/functions
then
   . %INITDIR%/functions
   OUTPUT_TYPE="2"
elif test -f /etc/rc.status
then
   . /etc/rc.status
   OUTPUT_TYPE="1"
elif test -f %INITDIR%/rc.status
then
   . %INITDIR%/rc.status
   OUTPUT_TYPE="1"
elif test -f /sbin/functions.sh
then
   . /sbin/functions.sh
   OUTPUT_TYPE="4"
elif test -f %INITDIR%/functions.sh
then
   . %INITDIR%/functions.sh
   OUTPUT_TYPE="4"
else
   OUTPUT_TYPE="0"
fi

#
# Helper functions to display results
#
OUTPUT_LABEL=""
Output_Message ()
{
   case "$OUTPUT_TYPE" in
      1)
         echo -n $1
      ;;
      2)
         echo -n $1
      ;;
      3)
         OUTPUT_LABEL="$1"
      ;;
      4)
         ebegin $1
      ;;
      *)
      ;;
   esac
}

Output_Result ()
{
   case "$OUTPUT_TYPE" in
      1)
         case "$RETVAL" in
            0)
               echo -e "$rc_done"
            ;;
            7)
               echo -e "$rc_failed"
            ;;
            *)
            ;;
         esac
      ;;
      2)
         case "$RETVAL" in
            0)
               echo_success
            ;;
            *)
               echo_failure
            ;;
         esac
         echo ""
      ;;
      3)
         case "$RETVAL" in
            0)
               log_success_msg $OUTPUT_LABEL
            ;;
            *)
               log_failure_msg $OUTPUT_LABEL
            ;;
         esac
      ;;
      4)
         eend $RETVAL
      ;;
      *)
      ;;
   esac
}

#
# Make sure everything is there
#
[ -d %BINDIR% ] || exit 5
[ -x %BINDIR%/TwDriver ] || exit 5
[ -x %BINDIR%/TWXinputInstall.perl ] || exit 5

#
# Do what we were asked to do
#
case "$1" in
   start)
      Output_Message "Starting TwDriver"
      %BINDIR%/TWXinputInstall.perl
      [ x%TWDRV_CFG% != x ] && [ -r %TWDRV_CFG% ] && . %TWDRV_CFG%
      [ -d /dev/shm ] && rm -f /dev/shm/*TwConfig*
      %BINDIR%/TwDriver $TWDRV_OPTIONS
      RETVAL=$?
      Output_Result
      ;;
   stop)
      Output_Message "Stopping TwDriver"
      %BINDIR%/TwDriver stop
      RETVAL=$?
      Output_Result
      ;;
   status)
      Output_Message "Checking for TwDriver"
      %BINDIR%/TwDriver check
      RETVAL=$?
      Output_Result
      ;;
   restart)
      $0 stop
      $0 start
      RETVAL=$?
      ;;
   try-restart)
      $0 status
      if test $? = 0
      then
         $0 restart
         RETVAL=$?
      else
         RETVAL=0
      fi
      ;;
   *)
      echo "Usage: $0 {start|stop|status|restart|try-restart}"
      exit 1
esac

exit $RETVAL
